# rpm/SOURCES/etc/cron.d/cron-logrotate-contextBroker-size
# cron job to test and execute log rotate by size.
*/30 * * * * orion logrotate  /etc/sysconfig/logrotate-contextBroker-size